# com.yanyao.ai-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yanyao.ai-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yanyao/Projects/side/ai-bridge/bin/ai-bridge</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yanyao/Projects/side/ai-bridge</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yanyao</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>AI_BRIDGE_SOCKS</key>
        <string>127.0.0.1:2333</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/yanyao/Projects/side/ai-bridge/logs/cpa.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yanyao/Projects/side/ai-bridge/logs/cpa.err.log</string>
</dict>
</plist>
